# deploy/systemd/python-domaintools.service
[Unit]
Description=Gunicorn for Python DomainTools
After=network.target redis.service
Wants=redis.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/python.domaintools.site/htdocs
Environment="PATH=/var/www/python.domaintools.site/htdocs/.venv/bin"
EnvironmentFile=/etc/domaintools/domaintools.env
ExecStart=/var/www/python.domaintools.site/htdocs/.venv/bin/gunicorn \
  --workers 3 \
  --bind 127.0.0.1:8000 \
  --timeout 120 \
  --graceful-timeout 30 \
  --max-requests 1000 \
  --max-requests-jitter 50 \
  wsgi:app
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target